# ryzenadj-gtk-apply.service
[Unit]
Description=Apply saved ryzenadj-gtk settings on boot
After=multi-user.target systemd-modules-load.service
Wants=systemd-modules-load.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
ExecStart=/usr/bin/python3 /usr/share/ryzenadj-gtk/ryzen.py --apply-saved

[Install]
WantedBy=multi-user.target
